# .config/docker/data/etc/systemd/system/gnome-keyring.service
[Unit]
Description=GNOME Keyring daemon (Secret Service provider)
After=dbus-user-session.service
Requires=dbus-user-session.service
Before=rekindle-init.service

[Service]
Type=forking
User=rekindle
Group=rekindle
Environment=HOME=/home/rekindle
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=/run/user/1000/bus
ExecStart=/usr/bin/gnome-keyring-daemon --start --components=secrets
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
