# .config/docker/data/etc/systemd/system/rekindle-node.service
[Unit]
Description=Rekindle node daemon (IPC bus + Veilid transport)
Documentation=https://github.com/scopecreep-zip/rekindle
After=network-online.target gnome-keyring.service dbus-user-session.service
Wants=network-online.target
Requires=gnome-keyring.service dbus-user-session.service
Before=rekindle-init.service
StartLimitBurst=5
StartLimitIntervalSec=300

[Service]
Type=notify
ExecStart=/usr/local/bin/rekindle node start
ExecStop=/usr/local/bin/rekindle node stop
Restart=on-failure
RestartSec=5
User=rekindle
Group=rekindle

# Startup: Veilid attach can take 30-60 seconds on slow networks.
# The daemon sends READY=1 in Locked state BEFORE transport starts.
# The unlock flow (triggered by init or user) handles transport start.
TimeoutStartSec=90
# Shutdown: graceful drain budget.
TimeoutStopSec=15
# Watchdog disabled: the unlock+Veilid attach flow can take 60+ seconds
# which would trigger a watchdog kill. The daemon is healthy in Locked
# state waiting for unlock — no watchdog needed.
# WatchdogSec=30

# Environment
Environment=HOME=/home/rekindle
Environment=RUST_LOG=rekindle=debug,rekindle_transport=debug,rekindle_node=debug,veilid_core=warn
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=/run/user/1000/bus
Environment=XDG_RUNTIME_DIR=/run/user/1000

# Security hardening
NoNewPrivileges=yes
LimitCORE=0
LimitNOFILE=65536
LimitMEMLOCK=64M
MemoryMax=512M

[Install]
WantedBy=multi-user.target
